# templates/startup/upstart.conf.erb
description "Chronos scheduler for Mesos"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

<% if @run_user -%>
setuid <%= @run_user %>
<% end-%>
<% if @run_group -%>
setgid <%= @run_group %>
<% end -%>

<%=
exec = "exec #{@launcher_path}"
exec += " --jar '#{@jar_file_path}'" if @jar_file_path
exec
%>
